# deploy/casperfinder-bot.service
[Unit]
Description=CasperFinder Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/casperfinder-bot
ExecStart=/opt/casperfinder-bot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
